# crontab/sample-in-case-of-no-fqdn-specified
# min hour day month week command
@reboot $HOME/nsupdateifchanged/bin/nsupdateifchanged -z your.zone.name your_interface_name
0-59/3 * * * * $HOME/nsupdateifchanged/bin/nsupdateifchanged -z your.zone.name your_interface_name
