# crontab/sample-in-case-of-no-fqdn-specified
# min hour day month week command
@reboot $HOME/nsupdateifchanged/bin/nsupdateifchanged -z your.zone.name your_interface_name
0,10,20,30,40,50 * * * * $HOME/nsupdateifchanged/bin/nsupdateifchanged -z your.zone.name your_interface_name
