# crontab/sample
# min hour day month week command
@reboot $HOME/nsupdateifchanged/bin/nsupdateifchanged
0-59/3 * * * * $HOME/nsupdateifchanged/bin/nsupdateifchanged
